# deploy/systemd/ocp-commerce-catalog-worker.service
[Unit]
Description=OCP Commerce Catalog Search Index Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/workspace/OCP-Catalog
EnvironmentFile=/home/ubuntu/workspace/OCP-Catalog/.env
ExecStart=/home/ubuntu/.bun/bin/bun run commerce:catalog:worker
Restart=always
RestartSec=5
NoNewPrivileges=true
LimitNOFILE=524288

[Install]
WantedBy=multi-user.target
